# docker/snp-prototype-mongodb/conf/bootstrap.conf
[program:bootstrap]
command=/usr/local/bin/mongo --nodb /bootstrap/data/test-data.js
; we use the "fast" shutdown signal SIGINT
stopsignal=INT
redirect_stderr=true
autorestart=false
startsecs=0
stdout_logfile=/etc/supervisor/conf.d/bootstrap-stdout.log
