# docker/conf/mongo.conf
[program:mongod]
command=/usr/bin/mongod
; we use the "fast" shutdown signal SIGINT
stopsignal=INT
redirect_stderr=true
autorestart=true
